# pypi/crontab
0 5 * * 1 bandersnatch mirror >> /var/log/pypi.log
# Don't remove this line .